# app/supervisor-app.conf
[program:app-uwsgi]
command=/usr/local/bin/uwsgi --ini /home/docker/code/uwsgi.ini
environment=LC_CTYPE=C.UTF-8

[program:nginx-app]
command=/usr/sbin/nginx
environment=LC_CTYPE=C.UTF-8
